# scripts/supernode-scripts/azure/wolkaz/wolk-startup-scripts/scripts/profile.d/nagios.sh
#pathmunge /usr/lib64/nagios/plugins/
pathmunge /usr/local/nagios/libexec
