# scripts/supernode-scripts/azure/wolkaz/wolk-startup-scripts/scripts/profile.d/pushcode.sh
pathmunge /var/www/vhosts/mdotm.com/scripts/utils
#pathmunge /var/www/vhosts/anand/mdotm.com/scripts/utils
#pathmunge /var/www/vhosts/mdotm.com/scripts/utils
