# scripts/systems/anand/utils/plasma/plasma
# chkconfig: 2345 20 80
# description: Start plasma http service
# place this file in /etc/init.d

# Source function library.
. /etc/init.d/functions

case "$1" in
start)
   /root/go/src/github.com/wolkdb/plasma/build/bin/plasma &
   echo $!>/var/run/plasma.pid
   ;;
stop)
   kill `cat /var/run/plasma.pid`
   rm /var/run/plasma.pid
   ;;
restart)
   $0 stop
   $0 start
   ;;
status)
   if [ -e /var/run/plasma.pid ]; then
      echo plasma is running, pid=`cat /var/run/plasma.pid`
   else
      echo plasma is NOT running
      exit 1
   fi
   ;;
*)
   echo "Usage: $0 {start|stop|status|restart}"
esac

exit 0
